# fs/etc/services.d/cadvisor/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid "$PUID:$PGID" /usr/bin/cadvisor $CADVISOR_ARGS